# termux_usr/out/data/data/com.termux/files/usr/var/service/ftpd/run
#!/data/data/com.termux/files/usr/bin/sh
exec busybox tcpsvd -vE 0.0.0.0 8021 busybox ftpd -w $HOME
